# package/debian/media-agent.postrm
#!/bin/sh
set -e

case "$1" in
    remove|purge)
        echo "[media-agent] postrm: $1"
        if command -v systemctl >/dev/null 2>&1; then
            systemctl disable media-agent.service || true
            systemctl daemon-reload || true
        fi
        ;;
    *)
        ;;
esac

exit 0
